# dist/vm/postinst
#!/bin/sh

systemctl enable raceintospace
systemctl start raceintospace

exit 0